# discord-vc-bot.service
[Unit]
Description=Discord VC Auto-Kick Bot
After=network.target

[Service]
Type=simple
User=discord-bot
Group=discord-bot
WorkingDirectory=/opt/discord-vc-bot
Environment="DISCORD_BOT_TOKEN=YOUR_BOT_TOKEN_HERE"
ExecStartPre=/opt/discord-vc-bot/venv/bin/pip install -r requirements.txt
ExecStart=/opt/discord-vc-bot/venv/bin/python /opt/discord-vc-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/discord-vc-bot/logs /opt/discord-vc-bot/guild_settings.json

[Install]
WantedBy=multi-user.target